# ufw/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/default/
-rw-rw-r--	root/root	etc/default/ufw
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/ufw
drwxr-xr-x	root/root	etc/ufw/
-rw-r-----	root/root	etc/ufw/after.init
-rw-r-----	root/root	etc/ufw/after.rules
-rw-r-----	root/root	etc/ufw/after6.rules
drwxr-xr-x	root/root	etc/ufw/applications.d/
-rw-rw-r--	root/root	etc/ufw/applications.d/ufw-bittorent
-rw-rw-r--	root/root	etc/ufw/applications.d/ufw-chat
-rw-rw-r--	root/root	etc/ufw/applications.d/ufw-directoryserver
-rw-rw-r--	root/root	etc/ufw/applications.d/ufw-dnsserver
-rw-rw-r--	root/root	etc/ufw/applications.d/ufw-fileserver
-rw-rw-r--	root/root	etc/ufw/applications.d/ufw-loginserver
-rw-rw-r--	root/root	etc/ufw/applications.d/ufw-mailserver
-rw-rw-r--	root/root	etc/ufw/applications.d/ufw-printserver
-rw-rw-r--	root/root	etc/ufw/applications.d/ufw-proxyserver
-rw-rw-r--	root/root	etc/ufw/applications.d/ufw-webserver
-rw-r-----	root/root	etc/ufw/before.init
-rw-r-----	root/root	etc/ufw/before.rules
-rw-r-----	root/root	etc/ufw/before6.rules
-rw-rw-r--	root/root	etc/ufw/sysctl.conf
-rw-rw-r--	root/root	etc/ufw/ufw.conf
-rw-r-----	root/root	etc/ufw/user.rules
-rw-r-----	root/root	etc/ufw/user6.rules
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/ufw/
-rwxrwxr-x	root/root	lib/ufw/ufw-init
-rwxrwxr-x	root/root	lib/ufw/ufw-init-functions
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/python3.10/
drwxr-xr-x	root/root	usr/lib/python3.10/site-packages/
drwxr-xr-x	root/root	usr/lib/python3.10/site-packages/ufw-0.36.2-py3.10.egg-info/
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw-0.36.2-py3.10.egg-info/PKG-INFO
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw-0.36.2-py3.10.egg-info/SOURCES.txt
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw-0.36.2-py3.10.egg-info/dependency_links.txt
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw-0.36.2-py3.10.egg-info/top_level.txt
drwxr-xr-x	root/root	usr/lib/python3.10/site-packages/ufw/
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/__init__.py (EMPTY)
drwxr-xr-x	root/root	usr/lib/python3.10/site-packages/ufw/__pycache__/
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/__pycache__/__init__.cpython-310.pyc
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/__pycache__/applications.cpython-310.pyc
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/__pycache__/backend.cpython-310.pyc
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/__pycache__/backend_iptables.cpython-310.pyc
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/__pycache__/common.cpython-310.pyc
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/__pycache__/frontend.cpython-310.pyc
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/__pycache__/parser.cpython-310.pyc
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/__pycache__/util.cpython-310.pyc
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/applications.py
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/backend.py
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/backend_iptables.py
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/common.py
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/frontend.py
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/parser.py
-rw-r--r--	root/root	usr/lib/python3.10/site-packages/ufw/util.py
drwxr-xr-x	root/root	usr/sbin/
-rwxrwxr-x	root/root	usr/sbin/ufw
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-rw-r--	root/root	usr/share/man/man8/ufw-framework.8.gz
-rw-rw-r--	root/root	usr/share/man/man8/ufw.8.gz
drwxr-xr-x	root/root	usr/share/ufw/
drwxr-xr-x	root/root	usr/share/ufw/iptables/
-rw-r-----	root/root	usr/share/ufw/iptables/after.rules
-rw-r-----	root/root	usr/share/ufw/iptables/after6.rules
-rw-r-----	root/root	usr/share/ufw/iptables/before.rules
-rw-r-----	root/root	usr/share/ufw/iptables/before6.rules
-rw-r-----	root/root	usr/share/ufw/iptables/user.rules
-rw-r-----	root/root	usr/share/ufw/iptables/user6.rules
